# rc.local
#!/bin/sh -e

/root/httprelay/start.sh &

exit 0